# linux/unifiedshield.service
[Unit]
Description=UnifiedShield Anti-Censorship VPN
Documentation=https://github.com/unifiedshield/unifiedshield
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
ExecStart=/usr/bin/unifiedshield --core xray
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=256

# Watchdog (auto-restart if daemon becomes unresponsive)
WatchdogSec=30

# Sandboxing
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=true
ProtectControlGroups=true
ReadWritePaths=/run/unifiedshield /etc/unifiedshield

# Runtime directory
RuntimeDirectory=unifiedshield
RuntimeDirectoryMode=0755

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unifiedshield

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target
